# crontab
PATH=/usr/local/bin:$PATH

0 23 * * * /Users/othree/os-projects/moztw-download-validation/go > /tmp/cron.log
